# deploy/northstar-planner-api.service
[Unit]
Description=NorthStar Planner API
After=network.target postgresql.service

[Service]
Type=simple
User=sysops
WorkingDirectory=/forge/northstar-planner
EnvironmentFile=/etc/northstar-planner/api.env
ExecStart=/usr/bin/node /forge/northstar-planner/server-dist/server/index.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
